# rootfs/etc/s6/.s6-svscan/finish
#!/bin/sh

echo "aaaa"
# unmount fuse
if [ 1 -eq `ps -ef | grep abc | grep "rclone mount" | wc -l` ]; then
  fusermount -u -z /rclone/remote
fi

