# rootfs/etc/services.d/.s6-svscan/finish
#!/usr/bin/with-contenv sh

umount -f /rclone/data
umount -f /rclone/remote

