# rootfs/etc/services.d/20-mount-mergerfs/finish
#!/usr/bin/with-contenv sh

# unmount fuse
if [ 1 -eq `ps -ef | grep abc | grep mergerfs | wc -l` ]; then
  fusermount -u -z /rclone/data
fi
