# rootfs/etc/services.d/20-mount-mergerfs/run
#!/usr/bin/with-contenv sh

mount_command="mergerfs -f -o async_read=false,use_ino,allow_other,auto_cache,func.getattr=newest,category.action=all,category.create=ff /rclone/local:/rclone/remote /rclone/data"

echo "Executing mount mergerfs => $mount_command"
exec s6-setuidgid abc $mount_command
